# deploy/macos/com.vaf.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vaf.orchestrator</string>
  <key>WorkingDirectory</key>
  <string>__VAF_ROOT__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>orchestrator.api:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>__VAF_PORT__</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>__VAF_ROOT__</string>
    <key>VAF_DB_PATH</key>
    <string>__VAF_ROOT__/db/agentflow.db</string>
    <key>VAF_RUNS_ROOT</key>
    <string>__VAF_ROOT__/data/runs</string>
    <key>VAF_MATERIAL_LIBRARY_ROOT</key>
    <string>__VAF_ROOT__/data/01_素材库/对标视频/manual_import</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/orchestrator.out.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/orchestrator.err.log</string>
</dict>
</plist>
